# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /etc/supervisor/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
